# services/systemd/my-service.service
[Unit]
Description=My Custom Service

[Service]
ExecStart=/path/to/your/script.sh
Restart=always

[Install]
WantedBy=multi-user.target
